# deploy/rn1-tracker.service
[Unit]
Description=RN1 Activity Monitor (Market Discovery)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polybot
WorkingDirectory=/home/polybot/polymarket-rn1
ExecStart=/home/polybot/polymarket-rn1/venv/bin/python run_rn1_tracker.py
Restart=always
RestartSec=5
Environment=PYTHONIOENCODING=utf-8
Environment=PYTHONUNBUFFERED=1

# Graceful shutdown
KillSignal=SIGTERM
TimeoutStopSec=10

# Logging
StandardOutput=append:/home/polybot/polymarket-rn1/rn1_tracker.log
StandardError=append:/home/polybot/polymarket-rn1/rn1_tracker.log

[Install]
WantedBy=multi-user.target
